# lowBattery.timer
[Unit]
Description=lowBattery

[Timer]
OnUnitActiveSec=1s
OnBootSec=1s

[Install]
WantedBy=timers.target
